# jer.cx.service
[Unit]
Description=Jer.cx HTTP service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/www/jer.cx

[Install]
WantedBy=multi-user.target
